# agent/deploy/Inventory-ScheduledTask.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  HardView — Inventar-Agent (GPO-Vorlage fuer Task Scheduler 1.2)

  Import:  schtasks /Create /TN "HardView\HardwareInventar" /XML Inventory-ScheduledTask.xml /RU SYSTEM
  oder per GPP: Computerkonfiguration > Einstellungen > Geplante Aufgaben > Neu > "Geplante Aufgabe (Windows Vista+)"
                und die Werte aus dieser Datei uebernehmen.

  VOR DEM IMPORT ANPASSEN (zwei Platzhalter in <Arguments>):
    1) Skriptquelle:  \\YOUR_DOMAIN.local\NETLOGON\HardView\Invoke-Inventory.ps1   (von SYSTEM/Domaenen-Computern lesbar)
    2) Inventar-Ziel: \\FILESERVER\Inventory$\incoming                         (UNC! KEIN G:\ — SYSTEM hat keine gemappten Laufwerke)

  Der Task laeuft als SYSTEM; der Share-Zugriff erfolgt als Computerkonto (DOMAENE\PC$),
  daher muss "Domaenen-Computer" Schreibrecht auf dem Ziel-Share haben (siehe README-Deployment.md).

  EMPFOHLEN: Skript signieren. Der SYSTEM-Lauf aus einer Freigabe macht die
  Skript-Integritaet sicherheitsrelevant (Schreibrecht auf die Skriptquelle
  strikt auf Domaenen-Admins beschraenken).
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>HardView IT-Asset-Management: woechentliche Hardware-Inventarisierung (nur Hardware-Metadaten, kein Nutzer-Tracking).</Description>
    <Author>HardView</Author>
    <URI>\HardView\HardwareInventar</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-01-04T12:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <RandomDelay>PT4H</RandomDelay>
      <ScheduleByWeek>
        <DaysOfWeek>
          <Sunday />
        </DaysOfWeek>
        <WeeksInterval>1</WeeksInterval>
      </ScheduleByWeek>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>true</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT1H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>%SystemRoot%\System32\WindowsPowerShell\v1.0\powershell.exe</Command>
      <Arguments>-NoProfile -NonInteractive -ExecutionPolicy AllSigned -WindowStyle Hidden -File "\\YOUR_DOMAIN.local\NETLOGON\HardView\Invoke-Inventory.ps1" -OutputDir "\\FILESERVER\Inventory$\incoming"</Arguments>
    </Exec>
  </Actions>
</Task>
